# php-secookie/debian/php-fkooman-secookie.cron.daily
#!/bin/sh
/usr/bin/find /var/lib/php/sessions -type f -name "sses_*" -mtime +0 -delete
